# sodola-exporter.service
[Unit]
Description=Sodola HTTP Exporter for Prometheus
Documentation=https://github.com/yurkowashere/sodola-exporter
After=network.target

[Service]
Type=simple
User=prometheus
Group=prometheus
ExecStart=/usr/bin/python3 /opt/sodola-exporter/sodola_http_exporter.py --host 0.0.0.0 --port 9118 --log-level INFO
WorkingDirectory=/opt/sodola-exporter
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
